# server/kvenno-backend.service
[Unit]
Description=Kvenno Backend API Server
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/kvenno-server
Environment="NODE_ENV=production"
Environment="HOME=/var/www"
EnvironmentFile=/opt/kvenno-server/.env
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kvenno-backend

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
ReadWritePaths=/tmp

[Install]
WantedBy=multi-user.target
